# ultra-tracker.service
[Unit]
Description=Ultra Tracker
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=simple
Environment="COMPOSE_PROJECT_DIR=/proj/ultra_tracker"
Environment="COMPOSE_FILE=docker-compose.yml"
ExecStart=/usr/bin/podman-compose up
ExecStop=/usr/bin/podman-compose down -t 2
Restart=always

[Install]
WantedBy=multi-user.target
